# vps/just-ship-updater.service
[Unit]
Description=Just Ship Update Agent
Documentation=https://github.com/yves-s/just-ship
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=claude-dev
Group=claude-dev
ExecStart=/usr/local/bin/just-ship-updater.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=just-ship-updater

# Environment
Environment=HOME=/home/claude-dev
WorkingDirectory=/home/claude-dev

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/claude-dev /tmp /usr/local/bin/just-ship-updater.sh
ProtectHome=no

[Install]
WantedBy=multi-user.target
