# vps/just-ship-bot.service
# =============================================================================
# just-ship-bot.service — systemd Unit for Telegram Bot
#
# Installation:
#   sudo cp vps/just-ship-bot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now just-ship-bot
#
# Voraussetzung:
#   /home/claude-dev/.env              # Globale Keys (ANTHROPIC_API_KEY, OPENAI_API_KEY)
#   /home/claude-dev/.env.telegram-bot # Bot-spezifische Config
#
# Bot-Config (.env.telegram-bot):
#   TELEGRAM_BOT_TOKEN=123456:ABC-DEF...
#   SUPABASE_URL=https://xxx.supabase.co
#   SUPABASE_SERVICE_KEY=eyJ...
#   BOARD_API_URL=https://app.just-ship.io
#   TELEGRAM_BOT_SECRET=your-shared-secret-from-board-settings
# =============================================================================

[Unit]
Description=Just Ship Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=claude-dev
Group=claude-dev
WorkingDirectory=/home/claude-dev/just-ship/apps/bot

# Globale Keys + bot-spezifische Config
EnvironmentFile=/home/claude-dev/.env
EnvironmentFile=/home/claude-dev/.env.telegram-bot

ExecStart=/usr/bin/npx tsx bot.ts

# Logging
StandardOutput=append:/home/claude-dev/pipeline-logs/telegram-bot.log
StandardError=append:/home/claude-dev/pipeline-logs/telegram-bot.log

# Automatischer Neustart
Restart=always
RestartSec=10s

# Sicherheit
NoNewPrivileges=true
PrivateTmp=true

# Ressourcen-Limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
